# scripts/launchd/ai.myopenclaw.collect-agentops.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.myopenclaw.collect-agentops</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd __MYOPENCLAW_DIR__ &amp;&amp; __PYTHON3__ scripts/collect_agentops.py &gt;&gt; __MYOPENCLAW_DIR__/logs/collect-agentops.log 2&gt;&amp;1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__MYOPENCLAW_DIR__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__MYOPENCLAW_DIR__/logs/collect-agentops.log</string>

    <key>StandardErrorPath</key>
    <string>__MYOPENCLAW_DIR__/logs/collect-agentops.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
